# cookbooks/etcd/templates/default/etcd_cluster.service.erb
[Unit]
Description=etcd
Documentation=https://github.com/coreos/etcd
Conflicts=etcd.service

[Service]
Type=notify
Restart=always
RestartSec=5s
LimitNOFILE=40000
TimeoutStartSec=0

ExecStart=<%= @command %>

[Install]
WantedBy=multi-user.target
